# agents-chat-v2.service
[Unit]
Description=Agent Chat Server v2 (ACV2)
After=network.target

[Service]
Type=simple
User=spencer
WorkingDirectory=/home/spencer/agents-chat-v2
ExecStart=/home/spencer/agents-chat-v2/venv/bin/python main_v2.py
Restart=always
RestartSec=5
StandardOutput=append:/home/spencer/agents-chat-v2/server.log
StandardError=append:/home/spencer/agents-chat-v2/server.log
# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
